# s6/debian-root/etc/services.d/lighttpd/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping lighttpd"
killall -9 lighttpd
