# s6/debian-root/etc/services.d/lighttpd/run
#!/usr/bin/with-contenv bash

s6-echo "Starting lighttpd"
lighttpd -D -f /etc/lighttpd/lighttpd.conf
