# s6/debian-root/etc/services.d/cloudflared/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping cloudflared proxy"
killall -9 cloudflared
