# s6/debian-root/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bash

s6-echo "Starting cloudflared dns proxy"
/sbin/cloudflared proxy-dns --port 5533
