# s6/debian-root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash

s6-echo "Starting dnsmasq"

s6-setuidgid root dnsmasq -7 /etc/dnsmasq.d --no-daemon
