# src/s6/debian-root/etc/s6-overlay/s6-rc.d/_gravityonboot/up
#!/command/execlineb
background { bash -e /usr/local/bin/gravityonboot.sh }