# src/s6/debian-root/etc/s6-overlay/s6-rc.d/cron/finish
#!/command/with-contenv bash

s6-echo "Stopping cron"
killall --signal 9 cron
