# src/s6/debian-root/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv bash
exec -c
fdmove -c 2 1 /usr/sbin/cron -f
