# s6/alpine-root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

s6-echo "Starting php-fpm"

php-fpm5 -d daemonize=no
