# s6/debian-root/etc/s6-overlay/s6-rc.d/my-init/up
#!/command/execlineb -P

/etc/s6-overlay/scripts/my-init-up