# s6/debian-root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash
s6-echo "Starting crond"

exec -c
fdmove -c 2 1 /usr/sbin/cron -f
