# src/s6/debian-root/etc/s6-overlay/s6-rc.d/lighttpd/finish
#!/command/with-contenv bash

s6-echo "Stopping lighttpd"
service lighttpd-access-log stop
service lighttpd-error-log stop
killall --signal 9 lighttpd
