# s6/debian-root/etc/services.d/pihole-FTL/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping pihole-FTL"
/etc/init.d/pihole-FTL stop
