# s6/debian-root/etc/services.d/pihole-FTL/run
#!/usr/bin/with-contenv bash

DNSMASQ_USER=${DNSMASQ_USER:-root}

s6-echo "Starting pihole-FTL ($FTL_CMD) as ${DNSMASQ_USER}"
# /var/log/pihole*.log has this data - this duplicates it in docker logs unless you send to dev null
s6-setuidgid ${DNSMASQ_USER} pihole-FTL $FTL_CMD >/dev/null 2>&1
