# alpine/root/etc/services.d/dnsmasq/finish
#!/usr/bin/with-contenv bash

kill -9 $(pgrep dnsmasq)
