# s6/debian-root/etc/s6-overlay/s6-rc.d/1-startup/up
#!/usr/bin/execlineb
foreground { bash -e /start.sh }