# src/s6/debian-root/etc/s6-overlay/s6-rc.d/_postFTL/up
#!/command/execlineb
background { bash -e /usr/local/bin/_postFTL.sh }